# dist/init/upstart/clusterware-configurator.conf
################################################################################
##
## Alces Clusterware - Upstart service unit file
## Copyright (c) 2015 Alces Software Ltd
##
################################################################################
# Alces Clusterware boot configurator (Upstart unit)
description "Alces Clusterware boot configurator"
start on (local-filesystems and net-device-up IFACE!=lo)
task

script
  [ -f /etc/xdg/clusterware/config.vars.sh ] && . /etc/xdg/clusterware/config.vars.sh
  "$cw_ROOT"/libexec/share/configurator
endscript
