# dist/init/upstart/clusterware-executor.conf
################################################################################
##
## Alces Clusterware - Upstart service unit file
## Copyright (c) 2015 Alces Software Ltd
##
################################################################################
# Clusterware executor (Upstart unit)
description "Execute Clusterware task(s)"
start on stopped clusterware-configurator
task

env PLUGIN_PATH=/opt/clusterware/etc/handlers
exec /opt/clusterware/opt/pluginhook/bin/pluginhook execute
