# dist/init/upstart/clusterware-configurator.conf
################################################################################
##
## Alces Clusterware - Upstart service unit file
## Copyright (c) 2015 Alces Software Ltd
##
################################################################################
# Clusterware configurator (Upstart unit)
description "Apply Clusterware configuration"
start on stopped cloud-final
task

env PLUGIN_PATH=/opt/clusterware/etc/handlers
exec /opt/clusterware/opt/pluginhook/bin/pluginhook configure
