# etc/profile.d/00-aa-disable-echo.sh
################################################################################
##
## Alces Clusterware - Shell configuration
## Copyright (c) 2018 Alces Software Ltd
##
################################################################################
stty -echo
