# etc/profile.d/06-pdsh.sh
################################################################################
##
## Alces Clusterware - Shell configuration
## Copyright (c) 2008-2015 Alces Software Ltd
##
################################################################################
export GENDERS_FILE=/opt/clusterware/etc/genders
export PDSH_GENDERS_FILE=$GENDERS_FILE
export PATH=/opt/clusterware/opt/pdsh/bin:/opt/clusterware/opt/genders/bin:$PATH
