# etc/profile.d/pdsh.csh
################################################################################
##
## Alces Clusterware - Shell configuration
## Copyright (c) 2008-2015 Alces Software Ltd
##
################################################################################
setenv GENDERS_FILE /opt/clusterware/etc/genders
setenv PDSH_GENDERS_FILE $GENDERS_FILE
setenv PATH /opt/clusterware/opt/pdsh/bin:/opt/clusterware/opt/genders/bin:$PATH
