# etc/profile.d/sshkey.sh
################################################################################
##
## Alces Clusterware - Shell configuration
## Copyright (c) 2015 Alces Software Ltd
##
################################################################################
if [ -x "${cw_ROOT}"/libexec/share/setup-sshkey ]; then
    "${cw_ROOT}"/libexec/share/setup-sshkey
fi
